# systemd/bedtime.service
[Unit]
Description=Bedtime reminder

[Service]
Type=oneshot
ExecStart=/home/jakub/devel/scripts/remind.sh

